# deploy/lociterm.service
[Unit]
Description=LociTerm - Web-based Multi-Terminal Server
After=network.target

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/lociterm --host 127.0.0.1 --port 8080 --data-dir /var/lib/lociterm
Restart=on-failure
RestartSec=5
Environment=HOME=%h
WorkingDirectory=%h

[Install]
WantedBy=multi-user.target
